# supervisor/beanstalkd.conf
[program:beanstalkd]
command=/usr/bin/beanstalkd -l 127.0.0.1 -p 11300 -u nobody
autostart=true
autorestart=true
